# profile.d/hadoop-env.sh
alias start-hs='mapred --daemon start historyserver';
alias stop-hs='mapred --daemon stop historyserver';
