# etc/profile.d/java-env.sh
export JAVA_HOME="/opt/jdk"
export JAVA_OPTS="-XX:+UseG1GC"
export PATH=${JAVA_HOME}:${PATH}